# packaging/gb-proxy.service
[Unit]
Description=GB-proxy legacy browser HTTP proxy
Documentation=https://github.com/salvogendut/GB-proxy
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
User=gb-proxy
Group=gb-proxy
Environment=PYTHONUNBUFFERED=1
Environment=GB_PROXY_CONFIG=/etc/gb-proxy/config.py
Environment=GB_PROXY_CACHE_DIR=/var/cache/gb-proxy
Environment=GB_PROXY_STATE_DIR=/var/lib/gb-proxy
EnvironmentFile=-/etc/sysconfig/gb-proxy
WorkingDirectory=/var/lib/gb-proxy
ExecStart=/usr/bin/gb-proxy
Restart=on-failure
RestartSec=5s
RestartPreventExitStatus=2

CacheDirectory=gb-proxy
CacheDirectoryMode=0750
StateDirectory=gb-proxy
StateDirectoryMode=0750
UMask=0077

# The proxy needs outbound network access but no elevated host privileges.
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=true
PrivateTmp=true
RemoveIPC=true
ProtectSystem=strict
ProtectHome=true
ProtectProc=invisible
ProcSubset=pid
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
